# tools/rdagent/com.omega.regime-book.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Daily regime book (shadow). Weekdays 09:00 pre-open. Proposes only.
  Install (after `bash regime_daily.sh` works by hand):
    cp com.omega.regime-book.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.omega.regime-book.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key><string>com.omega.regime-book</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jo/Omega/tools/rdagent/regime_daily.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
  </array>
  <key>StandardOutPath</key><string>/tmp/rda-regime.log</string>
  <key>StandardErrorPath</key><string>/tmp/rda-regime.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
